# banana-pi-setup/etc/rc.local
version https://git-lfs.github.com/spec/v1
oid sha256:2026bd4aa66f283bdefb12f9d11cc7b565831e829775e4dad0013896c6770cd6
size 547
