# raspberry-pi-setup/etc/rc.local
version https://git-lfs.github.com/spec/v1
oid sha256:d5059ca7181a301b73653bd3cac08e10399be6abe00a702b82cdbae6cfb4f89c
size 1051
